# lofi-pi-token-manager.service
[Unit]
Description=LoFi Pi Automatic Spotify Token Manager
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/lofi-pi
Environment=PYTHONPATH=/home/pi/lofi-pi
ExecStart=/usr/bin/python3 /home/pi/lofi-pi/auto_token_manager.py --continuous
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/lofi-pi

[Install]
WantedBy=multi-user.target
